# supervisord.conf
[supervisord]
nodaemon = true

[program:sshd]
command = /usr/sbin/sshd -D

[program:rsyslogd]
command = rsyslogd -n
